# _handover_20260614/systemd/vsp-netflow-api.service.d/10-env.conf
[Service]
EnvironmentFile=/home/devsecops/vsp-platform-v4.9.0/.env
ExecStart=
ExecStart=/usr/bin/gunicorn --workers 2 --bind 127.0.0.1:8930 vsp-netflow-api:app
