# _handover_20260614/systemd/vsp-netflow-api.service.d/override.conf
[Unit]
After=network-online.target postgresql@14-main.service
Wants=network-online.target

[Service]
Restart=on-failure
RestartSec=5
