# _handover_20260614/systemd/vsp-scanner-watchdog.service
[Unit]
Description=VSP Scanner Watchdog
[Service]
Type=oneshot
ExecStart=/usr/local/bin/vsp-scanner-watchdog.sh
