# _handover_20260614/systemd/vsp-scanner-watchdog.timer
[Unit]
Description=VSP Scanner Watchdog every 2 min
[Timer]
OnBootSec=60
OnUnitActiveSec=120
AccuracySec=15
[Install]
WantedBy=timers.target
