# _handover_20260614/systemd/vsp-audit-gate.service.d/10-jwtsecret.conf
[Service]
EnvironmentFile=/home/devsecops/vsp-platform-v4.9.0/.env
